# Software/rover_field_boot/deploy/rover-onboard.service
[Unit]
Description=Nisse onboard platform (status broker, OLED, audio, web UI)
After=network-online.target pulseaudio.service
Wants=network-online.target

[Service]
Type=simple
User=nisse
Group=nisse
Environment=NISSE_ROOT=/opt/nisse
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/nisse/venv/bin/rover-onboard \
    --rover-name ${ROVER_NAME} \
    --web-host 0.0.0.0 \
    --web-port 5000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
